# scripts/debian/dpkg-root/DEBIAN/postinst
#!/bin/sh
APP=blackflow

# Create directory
mkdir -p /var/blackflow/apps 2>/dev/null || :

# Install service
update-rc.d -f ${APP} remove
update-rc.d ${APP} defaults 30
service ${APP} start
