# scripts/debian/dpkg-root/DEBIAN/prerm
#!/bin/sh
APP=zwta

service ${APP} stop

if [ "$1" = "remove" ]; then
find /usr/lib/lyse_services/${APP} | grep "\.pyc$" | xargs rm
update-rc.d -f ${APP} remove
fi
